# screen-lock-monitor.service
[Unit]
Description=Monitor screen lock/unlock events
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/screen-lock-monitor.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
