# system-event-shutdown.service
[Unit]
Description=Log system shutdown event
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo "[$(date +%%Y-%%m-%%d\ %%H:%%M:%%S)] SYSTEM SHUTDOWN - Host: $(hostname)" >> /home/%u/login_logout_times.log'
RemainAfterExit=yes

[Install]
WantedBy=shutdown.target reboot.target halt.target
