# system-event-boot.service
[Unit]
Description=Log system boot event
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo "[$(date +%%Y-%%m-%%d\ %%H:%%M:%%S)] SYSTEM BOOT - Host: $(hostname)" >> /home/%u/login_logout_times.log'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
