# system-event-login.service
[Unit]
Description=Log user login event
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/system-event-logger.sh login
RemainAfterExit=yes

[Install]
WantedBy=default.target
