# system-event-logout.service
[Unit]
Description=Log user logout event
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/system-event-logger.sh logout
RemainAfterExit=yes

[Install]
WantedBy=default.target
